# metronome.service.in
[Unit]
Description=Metronome
Documentation=https://github.com/ahupowerdns/metronome/
After=network-online.target

[Service]
Type=notify
ExecStart=@bindir@/metronome --daemon 0 \
  --carbon-address '[::]:2003' \
  --webserver-address '[::]:8000' \
  --stats-directory /var/lib/metronome
Restart=on-failure
User=metronome
Group=metronome
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true
PrivateDevices=true
CapabilityBoundingSet=CAP_DAC_OVERRIDE CAP_IPC_LOCK CAP_KILL CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
